# scripts/agents-chat.service
[Unit]
Description=Agents-Chat (Next.js multi-agent chat UI)
Documentation=https://github.com/huanyingtianhe/agents-chat
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__GROUP__
WorkingDirectory=__PROJECT_DIR__

# Optional: drop overrides into /etc/agents-chat.env (KEY=VALUE per line)
EnvironmentFile=-/etc/agents-chat.env

# Sensible defaults (overridden by EnvironmentFile or by .env.local that Next loads)
Environment=NODE_ENV=production
Environment=PORT=3010
Environment=LOG_LEVEL=info
Environment=LOG_DIR=__PROJECT_DIR__/logs
Environment=LOG_FILE=app.log
Environment=LOG_ROTATE_FREQUENCY=daily
Environment=LOG_ROTATE_SIZE=10m
Environment=LOG_RETENTION=7

ExecStart=__SCRIPT_DIR__/start.sh --no-build

Restart=on-failure
RestartSec=10
TimeoutStopSec=30
KillSignal=SIGTERM

# journald captures stdout/stderr in addition to start.sh's tee'd files
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agents-chat

# Light sandboxing
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
